# .config/systemd/user/battery-notification.service
[Unit]
Description=Vérification batterie et notifications

[Service]
Type=oneshot
ExecStart=/home/benn/.config/mako/battery-notification.sh
